# packages/rollups-node/s6-overlay/s6-rc.d/redis/finish
#!/command/with-contenv sh
s6-setuidgid redis redis-cli shutdown