# packages/cartesi-node/s6-overlay/s6-rc.d/caddy/run
#!/command/with-contenv sh
exec \
    su nobody \
        -s /bin/sh \
        -c "caddy run --config /Caddyfile --adapter caddyfile 2>&1 | sed -e 's/^/caddy: /'"
