# packages/cartesi-node/s6-overlay/s6-rc.d/caddy/run
#!/command/with-contenv sh
exec s6-setuidgid cartesi caddy run --config /Caddyfile --adapter caddyfile
