# packages/cartesi-node/s6-overlay/s6-rc.d/graphql-server/run
#!/command/with-contenv sh
exec \
    su nobody \
        -s /bin/sh \
        -c "graphql-server 2>&1 | sed -e 's/^/graphql-sever: /'"
