# packages/cartesi-node/s6-overlay/s6-rc.d/graphql-server/run
#!/command/with-contenv sh
s6-echo "Starting rollups-graphql-server..."
exec s6-setuidgid cartesi graphql-server
