# packages/cartesi-node/s6-overlay/s6-rc.d/advance-runner/run
#!/command/with-contenv sh
s6-echo "Starting rollups-advance-runner..."
exec s6-setuidgid cartesi advance-runner --health-check-port=8083
