# packages/cartesi-node/s6-overlay/s6-rc.d/advance-runner/run
#!/command/with-contenv sh
exec \
    su nobody \
        -s /bin/sh \
        -c "advance-runner --health-check-port=8083 2>&1 | sed -e 's/^/advance-runner: /'"
