# packages/cartesi-node/s6-overlay/s6-rc.d/state-server/run
#!/command/with-contenv sh
exec \
    su nobody \
        -s /bin/sh \
        -c "rollups_state_server 2>&1 | sed -e 's/^/state-server: /'"