# packages/cartesi-node/s6-overlay/s6-rc.d/state-server/run
#!/command/with-contenv sh
s6-echo "Starting rollups-state-server..."
exec s6-setuidgid cartesi rollups_state_server
