# packages/rollups-node/s6-overlay/s6-rc.d/traefik/run
#!/command/with-contenv sh
exec s6-setuidgid nobody traefik --configfile /etc/traefik/traefik.yaml
