# packages/rollups-node/s6-overlay/s6-rc.d/advance-runner/run
#!/command/with-contenv sh
exec s6-notifyoncheck -w 500 -t 500 -n 5 \
    s6-setuidgid cartesi advance-runner --health-check-port=8083
