# packages/cartesi-node/s6-overlay/s6-rc.d/postgres/finish
#!/command/with-contenv sh

s6-echo "Stopping PostgreSQL server..."

exec /usr/bin/pg_ctlcluster 13 main stop -- -m fast