# packages/cartesi-node/s6-overlay/s6-rc.d/dispatcher/run
#!/command/with-contenv sh
s6-echo "Starting rollups-dispatcher..."
exec s6-setuidgid cartesi rollups_dispatcher --hc-port=8081
