# packages/cartesi-node/s6-overlay/s6-rc.d/dispatcher/run
#!/command/with-contenv sh
exec \
    su nobody \
        -s /bin/sh \
        -c "rollups_dispatcher --hc-port=8081 2>&1 | sed -e 's/^/rollups-dispatcher: /'"