# packages/rollups-node/s6-overlay/s6-rc.d/rollups-node/run
#!/command/with-contenv sh
set -e
exec s6-setuidgid cartesi cartesi-rollups-node
