# packages/cartesi-node/s6-overlay/s6-rc.d/redis/finish
#!/command/with-contenv sh

s6-echo "Stopping Redis server..."

exec su redis -s /bin/sh -c "redis-cli shutdown"
