# packages/cartesi-node/s6-overlay/s6-rc.d/redis/run
#!/command/with-contenv sh
s6-echo "Starting Redis server..."
cd /var/lib/redis
exec s6-setuidgid redis redis-server
