# packages/cartesi-node/s6-overlay/s6-rc.d/indexer/run
#!/command/with-contenv sh
s6-echo "Starting rollups-indexer..."
exec s6-setuidgid cartesi indexer --health-check-port=8082
