# packages/cartesi-node/s6-overlay/s6-rc.d/inspect-server/run
#!/command/with-contenv sh
s6-echo "Starting rollups-inspect-server..."
exec s6-setuidgid cartesi inspect-server
